# etc/supervisor/conf.d/apache.conf
[program:apache]
command=/scripts/apache/apache-start.sh

